# patroni.conf
[program:patroni]
command = patroni patroni.yaml
numprocs=1
directory = /var/lib/postgresql
autostart = true
user = postgres
stdout_logfile = /var/log/postgresql/patroni.log
#environment = 
